# applications/luci-app-advanced/root/etc/uci-defaults/80_advanced
#!/bin/sh

touch /etc/config/advanced
uci add advanced defaults
uci commit

return 0