# applications/luci-app-advanced/root/etc/uci-defaults/80_advanced
#!/bin/sh

CONFIG_FILE="/etc/config/advanced"

if [ ! -f "$CONFIG_FILE" ]; then

    touch "$CONFIG_FILE"
    
    uci add advanced defaults

    uci commit advanced
fi

if ! uci show advanced.@edcca[0] >/dev/null 2>&1; then

    uci add advanced edcca
    uci commit advanced
fi

if ! uci show advanced.@advanced_packet_steering[0] >/dev/null 2>&1; then
    
    uci add advanced advanced_packet_steering
    uci commit advanced
fi

return 0
